# backend/celery.conf
[program:celery]
directory = /app
command = celery -A tasks worker